# samizdat.rc
#!/bin/sh

# PROVIDE: samizdat
# REQUIRE: LOGIN postgresql redis
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable samizdat:
# samizdat_enable="YES"
# samizdat_user="www"
# samizdat_group="www"
# samizdat_dir="/usr/local/www/samizdat"
# samizdat_pidfile="/var/run/samizdat.pid"
# samizdat_logfile="/var/log/samizdat.log"

. /etc/rc.subr

name="samizdat"
rcvar=samizdat_enable

load_rc_config $name

: ${samizdat_enable="NO"}
: ${samizdat_user="www"}
: ${samizdat_group="www"}
: ${samizdat_dir="/sites/Samizdat"}
: ${samizdat_pidfile="${samizdat_dir}/bin/samizdat.pid"}
: ${samizdat_logfile="/var/log/samizdat.log"}

pidfile="${samizdat_pidfile}"
command="/usr/bin/su"
command_args="-m ${samizdat_user} /usr/local/bin/perl /usr/local/bin/hypnotoad ${samizdat_dir}/bin/samizdat"
procname="hypnotoad"
required_dirs="${samizdat_dir}"
samizdat_prestart()
{
    # Set environment for the daemon
    export MOJO_MODE=production
    export MOJO_HOME="${samizdat_dir}"
    export HOME="${samizdat_dir}"
    export PWD="${samizdat_dir}"
    
    # Ensure PID file directory exists
    install -d -o "${samizdat_user}" -g "${samizdat_group}" -m 755 "$(dirname ${samizdat_pidfile})"
    
    # Ensure log file exists with proper permissions  
    touch "${samizdat_logfile}"
    chown "${samizdat_user}:${samizdat_group}" "${samizdat_logfile}"
    chmod 640 "${samizdat_logfile}"
    
    # Change to app directory - this affects where su starts
    cd "${samizdat_dir}"
}

start_precmd="samizdat_prestart"

run_rc_command "$1"